# deploy.anongram.service
[Unit]
Description=Anongram Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/gram
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sh -c "/usr/bin/truncate -s 0 /root/gram/bot.log || true"
ExecStart=/root/gram/.venv/bin/python /root/gram/main.py
Restart=always
RestartSec=5
StandardOutput=append:/root/gram/bot.log
StandardError=append:/root/gram/bot.log

[Install]
WantedBy=multi-user.target
